# etc/cron/crontab
*/1 * * * * php /var/www/public/admin/cli/cron.php >/dev/null